# src/Shopper/WorkerService/.service
[Unit]
Description=my test app

[Service]
Type=notify
ExecStart=/usr/sbin/testapp

[Install]
WantedBy=multi-user.target